# setup/crontab.cron
# Pager start

@reboot pager/software/launch.sh

# Pager end
